# scripts/helpers/nvme.service
[Unit]
Description=NVME Service
Before=docker.service containerd.service

[Service]
ExecStart=/nvme.sh

[Install]
WantedBy=multi-user.target
